# ubuntu_study/autostart.conf
# auto startup script 
# put this file into '/etc/init'

start on startup
task
exec /path/to/script.py
